# launchd/com.darkcreation.gwn.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.darkcreation.gwn.hourly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/python3</string>
      <string>-m</string>
      <string>app.cli</string>
      <string>run-once</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONPATH</key>
      <string>/Users/darkcreation/Documents/german-word-notifier/src</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/darkcreation/Documents/german-word-notifier</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/gwn.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/gwn.stderr.log</string>
  </dict>
</plist>
